# systemd/user/dms-osd-unstick-patcher.service
[Unit]
Description=DMS OSD Unstick Patcher - fix "no OSD after hyprlock unlock" bug
Documentation=file://%h/.local/bin/dms-osd-unstick-patcher.sh
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dms-osd-unstick-patcher.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
