# systemd/user/dms.service
[Unit]
Description=Dank Material Shell (DMS)
StartLimitBurst=10
StartLimitIntervalSec=60

[Service]
Type=simple
Environment=QT_QPA_PLATFORMTHEME=gtk3
# Block until the Wayland socket exists, so we don't race Hyprland startup.
# This handles the timing problem that PartOf=graphical-session.target was
# supposed to solve — but graphical-session.target never auto-starts on this
# box (no uwsm), so binding to default.target + this socket-wait is the
# reliable pattern.
ExecStartPre=/bin/sh -c 'until [ -S "$XDG_RUNTIME_DIR/$WAYLAND_DISPLAY" ]; do sleep 0.2; done'
# Skip start if not under Hyprland: DMS depends on Hyprland IPC + protocols
# (foreign_toplevel, ext_workspace, ext_session_lock, wlr_gamma, etc). Under
# other compositors (e.g. y5/Nourish) it partly renders but the dock is empty
# and half the widgets do not update. XDG_CURRENT_DESKTOP is imported into the
# user systemd env by the compositor at session start, so it is populated by
# the time the WAYLAND_DISPLAY socket exists.
ExecStartPre=/bin/sh -c 'test "$XDG_CURRENT_DESKTOP" = "Hyprland" || exit 1'
ExecStart=/usr/bin/dms run
Restart=always
RestartSec=3
# Restart only kills the main dms process. Apps living in this cgroup
# (chromium, foot, terminals, claude session, etc) are left running so a
# DMS restart can pick up patched QML without nuking the whole session.
# Default KillMode=control-group would SIGTERM the entire tree.
KillMode=process

[Install]
WantedBy=default.target
